# backup.service
[Unit]
Description=Run backup

[Service]
ExecStart=/bin/bash /home/micah/scripts/backup.sh

[Install]
WantedBy=default.target
